# mac-remote.service
[Unit]
Description=mac-remote server (task queue for Claude -> Mac agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mywork/mytest/mac-remote
ExecStart=/opt/mywork/mytest/mac-remote/server.sh
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
